# patch_root/lib/systemd/system/usbtypec-otg.service
[Unit]
Description=USB Type C OTG helper script

[Service]
ExecStart=/usr/bin/typec_host_init.sh

[Install]
WantedBy=basic.target
